# files/logrotate.cron
#!/bin/bash

/usr/sbin/logrotate -m /usr/local/bin/maillog /etc/logrotate.conf

chown root:root /srv/save/logs/*
chmod 644 /srv/save/logs/*

if [ -d /tmp/maillogs ] ; then
  today=$(date +%Y%m%d)
  cd /tmp/maillogs/
    tar -czf /tmp/$(hostname)-$today.tar.gz *
  cd
  echo "Have some logs" | mail -s "Logs for $today" -r logs@$(hostname -f) -a /tmp/$(hostname)-$today.tar.gz magic@a-rwx.org
  rm -r /tmp/maillogs /tmp/$(hostname)-$today.tar.gz
fi

